# vendor/jeremy/v6.5.7/collector/ops/systemd/bdag-lan-p2p-firewall.service
[Unit]
Description=Allow LAN BlockDAG P2P sync ports
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
Environment=BDAG_LAN_CIDR=192.168.1.0/24
Environment=BDAG_P2P_PORTS=8151,8152
ExecStart=/usr/local/sbin/bdag-allow-lan-p2p-iptables
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
