# vendor/jeremy/v6.5.7/collector/ops/systemd/user-bdag-local-peers.timer
[Unit]
Description=Refresh BlockDAG local P2P peer discovery

[Timer]
OnBootSec=3min
OnUnitActiveSec=30min
Persistent=true
RandomizedDelaySec=5min
Unit=bdag-local-peers.service

[Install]
WantedBy=timers.target
