# vendor/jeremy/v6.3.20/ops/systemd/bdag-dashboard.service
[Unit]
Description=BlockDAG pool operations dashboard
Wants=network-online.target docker.service bdag-status-sampler.service
After=network-online.target docker.service bdag-status-sampler.service

[Service]
Type=simple
WorkingDirectory=/home/jeremy/blockdag-asic-pool
Environment=BDAG_PROJECT_ROOT=/home/jeremy/blockdag-asic-pool
Environment=BDAG_RUNTIME_DIR=/home/jeremy/blockdag-asic-pool/ops/runtime
Environment=BDAG_NODE_SERVICES=node
Environment=BDAG_STACK_SERVICES=postgres,node,pool
Environment=BDAG_POOL_CONTAINER=pool
Environment=BDAG_POOL_DB_CONTAINER=postgres
Environment=BDAG_DASHBOARD_BIND=127.0.0.1
Environment=BDAG_DASHBOARD_PORT=8088
EnvironmentFile=-/home/jeremy/blockdag-asic-pool/.env
EnvironmentFile=-/home/jeremy/blockdag-asic-pool/ops/runtime/ops.env
ExecStart=/usr/bin/env python3 /home/jeremy/blockdag-asic-pool/ops/dashboard.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
