# vendor/jeremy/v6.3.20/ops/systemd/user-bdag-chain-state-self-heal.service
[Unit]
Description=BlockDAG chain-state self-heal restore
After=default.target docker.service

[Service]
Type=oneshot
WorkingDirectory=%h/blockdag-asic-pool
Environment=BDAG_PROJECT_ROOT=%h/blockdag-asic-pool
Environment=BDAG_RUNTIME_DIR=%h/blockdag-asic-pool/ops/runtime
EnvironmentFile=-%h/blockdag-asic-pool/ops/runtime/ops.env
TimeoutStartSec=12h
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
CPUWeight=50
IOWeight=50
ExecStartPre=/bin/sh -c 'i=0; while [ "$i" -lt 60 ]; do (docker info >/dev/null 2>&1 || sudo -n docker info >/dev/null 2>&1) && exit 0; i=$((i+1)); sleep 5; done; exit 1'
ExecStart=%h/blockdag-asic-pool/ops/chain-state-self-heal.sh --from-systemd
