# vendor/jeremy/v6.3.20/ops/systemd/user-bdag-chain-restore-guard.timer
[Unit]
Description=Run BlockDAG chain restore-point freshness guard

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
RandomizedDelaySec=20s
Unit=bdag-chain-restore-guard.service

[Install]
WantedBy=timers.target
