# vendor/jeremy/v6.3.20/ops/systemd/user-bdag-rawdatadir-sidecar.timer
[Unit]
Description=Refresh BlockDAG raw datadir sidecar copy every 5 minutes

[Timer]
OnBootSec=5m
OnActiveSec=5m
OnUnitActiveSec=5m
RandomizedDelaySec=2m
Persistent=true
Unit=bdag-rawdatadir-sidecar.service

[Install]
WantedBy=timers.target
