# vendor/jeremy/v6.3.20/ops/systemd/user-bdag-stack-sentinel.timer
[Unit]
Description=Run BlockDAG mining stack sentinel

[Timer]
OnBootSec=2min
OnUnitActiveSec=60s
AccuracySec=15s
Persistent=true
RandomizedDelaySec=5s
Unit=bdag-stack-sentinel.service

[Install]
WantedBy=timers.target
