# vendor/jeremy/v6.3.20/collector/ops/systemd/user-bdag-sync-coordinator.timer
[Unit]
Description=Run BlockDAG sync coordinator during large catch-up

[Timer]
OnBootSec=3min
OnUnitActiveSec=2min
AccuracySec=30s
Persistent=true
Unit=bdag-sync-coordinator.service

[Install]
WantedBy=timers.target

