# vendor/jeremy/v6.3.20/collector/ops/systemd/user-codex-memory.service
[Unit]
Description=Codex context memory service
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/jeremy/blockdag-asic-pool
Environment=CODEX_PROJECT_ROOT=/home/jeremy/blockdag-asic-pool
Environment=CODEX_MEMORY_DB=/home/jeremy/.codex/memories/context-store/context.sqlite
Environment=CODEX_HISTORY_FILE=/home/jeremy/.codex/history.jsonl
Environment=CODEX_MEMORY_INTERVAL=10
ExecStart=/usr/bin/env python3 /home/jeremy/blockdag-asic-pool/ops/codex_memory.py watch --interval 10
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
