# vendor/jeremy/v6.5.7/ops/systemd/bdag-mining-host-tuning.timer
[Unit]
Description=Reapply passive BlockDAG mining host latency tuning

[Timer]
OnBootSec=90s
OnUnitActiveSec=5min
AccuracySec=30s
RandomizedDelaySec=30s
Unit=bdag-mining-host-tuning.service

[Install]
WantedBy=timers.target
