# vendor/jeremy/v6.5.7/ops/systemd/user-bdag-hourly-snapshot.timer
[Unit]
Description=Run BlockDAG chain snapshot periodically

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=120
Unit=bdag-hourly-snapshot.service

[Install]
WantedBy=timers.target
