# vendor/jeremy/v6.5.7/ops/systemd/user-bdag-ipfs-content-sidecar.timer
[Unit]
Description=Run BlockDAG lazy IPFS content sidecar

[Timer]
OnBootSec=20m
OnActiveSec=30m
OnUnitActiveSec=30m
RandomizedDelaySec=10m
Persistent=false
Unit=bdag-ipfs-content-sidecar.service

[Install]
WantedBy=timers.target
