# vendor/jeremy/v6.5.7/ops/systemd/user-bdag-rawdatadir-source.service
[Unit]
Description=BlockDAG raw datadir FastArtifact source publisher
After=default.target docker.service

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/bdag-rawdatadir-sidecar.env
Nice=19
IOSchedulingClass=idle
CPUWeight=10
IOWeight=10
MemoryHigh=2G
TimeoutStartSec=4h
ExecStart=/usr/bin/env bash -lc 'cd "${BDAG_PROJECT_ROOT:?set BDAG_PROJECT_ROOT in ~/.config/bdag-rawdatadir-sidecar.env}" && exec ./ops/publish-rawdatadir-artifact.sh'
