# vendor/jeremy/v6.5.7/ops/systemd/user-bdag-mining-30min-guard.timer
[Unit]
Description=Run BlockDAG thirty-minute mining proof-of-health guard

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true
RandomizedDelaySec=2m
Unit=bdag-mining-30min-guard.service

[Install]
WantedBy=timers.target
